# webhook.service.example
[Unit]
Description=Mirror GHE Repositories

[Service]
ExecStart=/usr/bin/npm start
WorkingDirectory=
Environment=COMMAND=./mirror.sh "SECRET=YOUR SECRET HERE" PORT=8081
